# scripts/com.dotfiles.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Dotfiles auto-sync launchd agent.

  Installed to ~/Library/LaunchAgents/ by install.sh.
  The DOTFILES_DIR and HOME values are templated at install time
  (see install.sh Module 15: Auto-Sync).
-->
<dict>
    <key>Label</key>
    <string>com.dotfiles.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>__DOTFILES_DIR__/scripts/auto-sync.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DOTFILES_DIR</key>
        <string>__DOTFILES_DIR__</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SSH_AUTH_SOCK</key>
        <string>__HOME__/Library/Group Containers/2BUA8C4S2C.com.1password/t/agent.sock</string>
    </dict>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/share/dotfiles/sync.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/share/dotfiles/sync.log</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
